# supervisord/supervisord.d/overlord.ini
[group:overlord]
programs=overlord-watch-projects,overlord-watch-vm,overlord-poll-jails,overlord-poll-jail-info,overlord-poll-jail-stats,overlord-poll-jail-extras,overlord-poll-projects,overlord-poll-project-info,overlord-poll-autoscale,overlord-serve

[program:overlord-watch-projects]
command=/usr/local/bin/overlord watch-projects
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true

[program:overlord-watch-vm]
command=/usr/local/bin/overlord watch-vm
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true

[program:overlord-poll-jails]
command=/usr/local/bin/overlord poll-jails
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true

[program:overlord-poll-jail-info]
command=/usr/local/bin/overlord poll-jail-info
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true

[program:overlord-poll-jail-stats]
command=/usr/local/bin/overlord poll-jail-stats
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true

[program:overlord-poll-jail-extras]
command=/usr/local/bin/overlord poll-jail-extras
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true

[program:overlord-poll-projects]
command=/usr/local/bin/overlord poll-projects
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true

[program:overlord-poll-project-info]
command=/usr/local/bin/overlord poll-project-info
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true

[program:overlord-poll-autoscale]
command=/usr/local/bin/overlord poll-autoscale
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true

[program:overlord-serve]
command=/usr/local/bin/overlord serve
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/overlord/%(program_name)s.log
environment=OVERLORD_CONFIG=/usr/local/etc/overlord.yml,PATH="%(ENV_PATH)s:/usr/local/bin",HOME=/root,USER=root
stopasgroup=true
